# deploy/tasks.service
# systemd unit for the Tasks + Calendar app. Install to /etc/systemd/system/.
# Binds 127.0.0.1:8080 (see tasksd/__main__.py); public exposure is via the
# cloudflared tunnel -> public Caddy site (see Caddyfile.snippet), never direct.
[Unit]
Description=Tasks + Calendar web app (tasksd)
After=network-online.target radicale.service
Wants=network-online.target

[Service]
Type=exec
User=nicholaskmitchell
Group=nicholaskmitchell
WorkingDirectory=/home/nicholaskmitchell/tasks/backend
EnvironmentFile=/etc/tasks/tasks.env
ExecStart=/home/nicholaskmitchell/tasks/backend/.venv/bin/python -m tasksd
Restart=on-failure
RestartSec=3

# Hardening. ProtectHome=read-only keeps the app from writing anywhere under
# /home except the one path it needs (the SQLite cache). The app parses
# attacker-influenced iCalendar/HTTP and is internet-reachable through the
# tunnel, so the sandbox is tight: no capabilities (it binds 127.0.0.1:8080,
# nothing privileged), no devices, no namespaces, loopback only.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/nicholaskmitchell/tasks/backend
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
ProtectProc=invisible
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
CapabilityBoundingSet=
UMask=0077
# Loopback only: Caddy in, Radicale out. If TASKS_ACCESS_REQUIRED is ever
# turned on, DELETE these two lines — the Cloudflare Access JWKS fetch needs
# outbound HTTPS.
IPAddressAllow=localhost
IPAddressDeny=any

[Install]
WantedBy=multi-user.target
